# pi4/scripts/ogle-bridge.service
[Unit]
# OGLE vision node -> eye gaze bridge (RD-033). Reads /dev/ttyIRIS_OGLE, forwards
# GAZE:x,y / GAZE:LOST to the assistant CMD listener (UDP 127.0.0.1:10500).
# Requires assistant.service (owns the eyes port + runs the CMD listener).
Description=IRIS OGLE gaze bridge
After=assistant.service
Wants=assistant.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/python3 /home/pi/hardware/ogle_bridge.py
Restart=on-failure
RestartSec=3
# Tuning knobs — active values reflect bench findings (S137):
Environment=OGLE_FACING_REQUIRED=0
Environment=OGLE_LOST_TIMEOUT_S=2.5
# Environment=OGLE_CONF_GATE=60
# Environment=OGLE_FLIP_X=1
# Environment=OGLE_Y_BIAS=-0.10
# Environment=IRIS_OGLE_DEBUG=1

[Install]
WantedBy=multi-user.target
